# .launchd/com.newsweek.confiant-blocklist.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.newsweek.confiant-blocklist</string>

    <!-- Daily at 04:00 local time. Confiant's reports refresh after midnight
         UTC; running at 4 AM ET (08:00 UTC) means we have a fresh 24-hour
         window of data + blocks land before the work-day traffic peak.
         Originally weekly Mondays; moved to daily 2026-06-04 to shorten
         time-to-block on Google Cloaked + Phishing creatives.
         Adjust Hour/Minute/Weekday to taste; add a Weekday key to revert
         to weekly. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- If the Mac was asleep at the trigger time, run when it wakes. -->
    <key>StartCalendarIntervalRunAtLoad</key>
    <false/>

    <!-- IMPORTANT: edit these paths and env vars before `launchctl load`ing.
         The plist intentionally uses absolute paths so launchd doesn't depend
         on shell startup files. -->

    <key>WorkingDirectory</key>
    <string>/Users/roger/code/yield-dashboard</string>

    <!-- Protection 28044902 is "Everything" — the production catch-all
         Protection that covers all Newsweek inventory. Verify in
         GAM > Delivery > Protections before changing. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>python3</string>
        <string>/Users/roger/code/yield-dashboard/confiant_blocklist.py</string>
        <string>--protection-id</string>
        <string>28044902</string>
        <string>--protection-label</string>
        <string>Everything</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- launchd does not source ~/.zshrc. Credentials below come from
             ~/code/yield-dashboard/.env via _load_dotenv() in the script:
                 AGENTMAIL_API_KEY, AGENTMAIL_INBOX_ID, GAM_NETWORK_ID,
                 CONFIANT_API_KEY, CONFIANT_REPORT_TO_EMAIL
             Do NOT redeclare them here as empty strings — setdefault won't
             override empty values and the post-run summary email will
             silently skip (which is how the daily summary went un-sent
             from launch until 2026-06-08). -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>CONFIANT_BLOCKLIST_PROFILE_DIR</key>
        <string>/Users/roger/.confiant-blocklist/playwright-profile</string>
        <key>CONFIANT_BLOCKLIST_STATE</key>
        <string>/Users/roger/.confiant-blocklist/state.sqlite</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/roger/.confiant-blocklist/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/roger/.confiant-blocklist/launchd.err.log</string>

    <!-- Keep launchd from killing a slow run; Playwright can take a minute. -->
    <key>AbandonProcessGroup</key>
    <true/>
</dict>
</plist>
